# deploy/supervisor/accel-queue.conf
; Скопировать в /etc/supervisor/conf.d/ и выполнить: supervisorctl reread && supervisorctl update
; Очередь whatsapp обязательна — иначе не обрабатываются SendOutboundMessageJob, ProcessWhatsappInboundJob и ответы после отклонения звонка.
[program:accel-queue-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/accel.kz/artisan queue:work database --queue=whatsapp,default --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/www/accel.kz/storage/logs/queue-worker.log
stopwaitsecs=3600
