# deploy/supervisor/accel-queue-provisioning.conf
; Дополнительный воркер только для очереди provisioning (SSL/nginx при создании тенанта).
; Скопировать в /etc/supervisor/conf.d/ вместе с accel-queue.conf
[program:accel-queue-provisioning]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/accel.kz/artisan queue:work database --queue=provisioning --sleep=2 --tries=2 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/accel.kz/storage/logs/queue-provisioning.log
stopwaitsecs=3600
